# DEBIAN/postinst
systemctl daemon-reload
systemctl enable dockerfilter-persistent