# DEBIAN/prerm
service dockerfilter-persistent stop
systemctl disable dockerfilter-persistent